# 90-tablet-config.rules
ACTION=="add", SUBSYSTEM=="usb", ATTR{idVendor}=="1f3a", ATTR{idProduct}=="1002", RUN+="/opt/ubislate_config/configure.sh"